# templates/default/upstart.conf.erb
#!upstart
description "<%= @project_name %>"

start on runlevel [2345]
start on (started networking)

stop on runlevel [016]
stop on (stopping networking)

respawn
respawn limit 0 60

script
    exec <%= @start_command %>
end script
