# templates/default/upstart.conf.erb
#!upstart
description "<%= @project_name %>"

start on startup
stop on shutdown

respawn

exec <%= @start_command %>